# supervisor/conf.d/dovecot.ini
[program:dovecot]
command=/usr/sbin/dovecot -F
process_name=%(program_name)s
numprocs=1
priority=999
startsecs=5
redirect_stderr=true
autostart=true
