# supervisor/conf.d/postfix.ini
[program:postfix]
command=/usr/sbin/postfix start
process_name=%(program_name)s
numprocs=1
priority=99
startsecs=1
redirect_stderr=true
autorestart=false
startretries=0
